# scripts/upstart.conf
description     "${upstart_description}"

# no start option as you might not want it to auto-start
# This might not be supported - you might need a: start on runlevel [3]
stop on runlevel [!2345]

start on startup
stop on shutdown

script
    cd ${upstart_path}
    ${upstart_command}
end script

# if you want it to automatically restart if it crashes, leave the next line in
respawn